# gateway/deploy/local-tester-gateway.service
[Unit]
Description=local-tester LLM gateway
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/local-tester-gateway.env
ExecStart=/usr/bin/node /opt/local-tester-gateway/dist/index.js
Restart=always
RestartSec=2
User=gateway
Group=gateway
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
